# fakewww_bind_dnsa_1.0.0/supervisord.conf
[supervisord]
nodaemon = true
logfile = /var/log/supervisord.log
pidfile = /run/supervisord.pid

[program:sshd]
command = /usr/sbin/sshd -D

[program:named]
command = /usr/sbin/named -u named

[program:rsyslog]
command = /usr/sbin/rsyslogd -n -f /etc/rsyslog.conf
